# systemd/.config/systemd/user/wallpaper-rotate.timer
[Unit]
Description=Rotate wallpaper every 2 hours

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target

